# Docker_emulator/conf.d/adb-utils.conf
[program:adb-utils]
command=./src/utils.sh
autorestart=false
#stdout_logfile=${LOG_PATH}/adb-utils.stdout.log
#stderr_logfile=${LOG_PATH}/adb-utils.stderr.log
priority=5