# Docker_emulator/conf.d/auto-recording.conf
[program:auto-recording]
command=./src/record.sh auto_record
autorestart=false
#stdout_logfile=${LOG_PATH}/video-recording.stdout.log
#stderr_logfile=${LOG_PATH}/video-recording.stderr.log
priority=4