# Docker_emulator/conf.d/docker-appium.conf
[program:docker-appium]
command=./src/appium.sh
autorestart=false
#stdout_logfile=${LOG_PATH}/docker-android.stdout.log
#stderr_logfile=${LOG_PATH}/docker-android.stderr.log
priority=4