# Docker_emulator/conf.d/openbox.conf
[program:openbox]
command=openbox-session
#stdout_logfile=${LOG_PATH}/openbox.stdout.log
#stderr_logfile=${LOG_PATH}/openbox.stderr.log
priority=2