# Docker_emulator/conf.d/port-forward.conf
[program:port-forward]
command=./src/port_forward.sh
autorestart=false
priority=1
