# Docker_emulator/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -display ${DISPLAY} -nopw -forever -shared
#stdout_logfile=${LOG_PATH}/x11vnc.stdout.log
#stderr_logfile=${LOG_PATH}/x11vnc.stderr.log
priority=2